# app/scripts/app_job.cron
# Supercronic job schedule
0 */6 * * * /app/scripts/fetch_and_publish.sh
