# app/scripts/app_job.cron
# Supercronic job schedule
*/1 * * * * /app/scripts/fetch_and_publish.sh
